# files/system_files/raptor-powerprofile.service
[Unit]
Description=Raptor OS Initial Power Profile
# Set a sane baseline power profile before Cortex or the user changes it.
# Without this, power-profiles-daemon starts in whatever state it was left in.
After=power-profiles-daemon.service
Requires=power-profiles-daemon.service

[Service]
Type=oneshot
# Set balanced as the boot default. Cortex overrides this at runtime.
ExecStart=/bin/bash -c 'powerprofilesctl set balanced 2>/dev/null || true'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
