# files/system_files/raptor-gpu-profile.service
[Unit]
Description=Raptor OS — GPU Profile Detection & Configuration
After=sysinit.target
Before=display-manager.service

[Service]
Type=oneshot
# Script is at /usr/lib/raptor/gpu-detect.sh — written by raptor-gpu-profile.sh
# at build time. /usr/bin/raptor-gpu-profile.sh is the old path (no longer used).
ExecStart=/usr/lib/raptor/gpu-detect.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
