# files/system_files/raptor-ensure-services.service
[Unit]
Description=Raptor OS Auto-Enable Services
After=local-fs.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'systemctl enable raptor-gpu-profile.service || logger -t raptor-ensure-services -p user.err "Failed to enable raptor-gpu-profile.service"'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
