# files/system_files/raptor-firstboot.service
[Unit]
Description=Raptor OS First Boot - Browser Choice
After=graphical.target
Wants=graphical.target
ConditionPathExists=!/var/lib/raptor-browser-choice-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/raptor-browser-choice.sh
ExecStartPost=/usr/bin/touch /var/lib/raptor-browser-choice-done
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
