# files/system_files/raptor-cpugovernor.service
[Unit]
Description=Raptor OS CPU Governor Management
After=sysinit.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/lib/raptor/gpu-detect.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
