# files/system_files/raptor-ensure-services.service
[Unit]
Description=Raptor OS Auto-Enable Services
After=local-fs.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
  for svc in raptor-gpu-profile.service raptor-powerprofile.service raptor-cpugovernor.service; do \
    if ! systemctl is-enabled "$svc" &>/dev/null; then \
      logger -t raptor-ensure-services "Enabling $svc"; \
      systemctl enable --now "$svc" || logger -t raptor-ensure-services -p user.err "Failed to enable $svc"; \
    fi; \
  done'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
