# system_files/raptor-firstboot.service
```ini
[Unit]
Description=Raptor OS First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/raptor-firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'rpm-ostree rebase ostree-unverified-registry:ghcr.io/cerberus9-dev/raptor-os:latest && touch /var/lib/raptor-firstboot-done && systemctl reboot'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
```
