# etc/systemd/user/pipewire.service.d/raptor-restart.conf
# Raptor OS: auto-restart PipeWire on failure (common after suspend/resume)
[Service]
Restart=on-failure
RestartSec=2
